# conf/systemd.service
[Unit]
Description=Shiori: simple bookmark manager
Requires=network-online.target
After=network-online.target

[Service]
Type=simple
User=__APP__
Group=__APP__
EnvironmentFile="__INSTALL_DIR__/.env"
WorkingDirectory=__INSTALL_DIR__/
ExecStart=__INSTALL_DIR__/shiori server
Restart=always

[Install]
WantedBy=multi-user.target
